# config/systemd/user/hyprlock.service


[Unit]
Description=User suspend actions
After=suspend.target

[Service]
Type=simple
ExecStart=/usr/bin/hyprlock
TimeoutSec=infinity
ExecStartPost=/usr/bin/sleep 1

[Install]
WantedBy=suspend.target
