# config/systemd/user/suspend.target
[Unit]
Description=User level suspend target
StopWhenUnneeded=yes
Wants=i3lock.service