# config/systemd/user/xorg.target
[Unit]
Description=is xorg running?
After=default.target
Requires=default.target

