# crontab/gwagner
0 2 * * * /bin/bash /home/gwagner/.bin/backup.sh
